# whisper-vox.service
[Unit]
Description=whisper-vox always-listening voice daemon
After=pipewire.service

[Service]
Type=simple
ExecStart=/usr/local/bin/wv
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
